# .config/service/ssh-agent/run
#!/bin/sh

exec ssh-agent -D -a "$XDG_RUNTIME_DIR/ssh-agent.socket"
