# services/keybase-redirector/run
#!/bin/sh

exec keybase-redirector /keybase
